# ApacheHDFS/wbench/Ubuntu-18.04+Hadoop-3.2.1+Kerberos-5/home/hadoop/bin/hdfs_namenode_simple.service
[Unit]
Description=HDFS namenode simple
Requires=network.target remote-fs.target
After=network.target remote-fs.target

[Service]
Type=forking
User=hadoop
Group=sudo
Environment=JAVA_HOME=/usr/
Environment=HADOOP_HOME=/home/hadoop/hadoop-3.2.1/
WorkingDirectory=/home/hadoop/
ExecStart=/home/hadoop/bin/start_hdfs_namenode_simple.sh
ExecStop=/home/hadoop/bin/stop_hdfs_namenode_simple.sh

[Install]
WantedBy=multi-user.target
Alias=hdfs_namenode_simple
